# supervisor/1/error.conf
[program:error]
command=/home/delkov/Desktop/test.sh
autostart=true
autorestart=true
stderr_logfile=/home/delkov/Desktop/ERROR_log_error.txt
stdout_logfile=/home/delkov/Desktop/ERROR_log.txt
