# supervisor/1/test.conf
[program:long_script]
command=/home/delkov/Desktop/test.sh
autostart=true
autorestart=true
stderr_logfile=/home/delkov/Desktop/log.txt
stdout_logfile=/home/delkov/Desktop/log_error.txt
